# scripts/launchd/com.assetflow.push-today-insights.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  launchd job: send a "today's insights" FCM push every morning at 08:00 KST.

  Install (one-time):
    bash scripts/launchd/install.sh

  Inspect:
    launchctl list | grep com.assetflow
    launchctl print gui/$(id -u)/com.assetflow.push-today-insights
    launchctl unload ~/Library/LaunchAgents/com.assetflow.push-today-insights.plist

  Notes:
    Depends on the production Next.js server listening on 127.0.0.1:3500.
    If the server is down at 08:00, the curl returns connection-refused —
    no retry. We treat this as best-effort.
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.assetflow.push-today-insights</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>__PROJECT_DIR__/scripts/push-today-insights.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>__PROJECT_DIR__</string>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>8</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>__PROJECT_DIR__/data/push-today-insights.log</string>
  <key>StandardErrorPath</key>
  <string>__PROJECT_DIR__/data/push-today-insights.log</string>

  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
