# scripts/launchd/com.assetflow.web.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  launchd job for the Asset Flow Next.js production server.

  Lifecycle:
    - RunAtLoad        : starts on user login / after install
    - KeepAlive        : auto-restarts if the process exits or crashes
    - ThrottleInterval : 10s minimum between restart attempts
                         (otherwise a crash-loop hammers the box)

  Listening on 127.0.0.1:3500, fronted by the Cloudflare Tunnel that
  maps assetflow.elkavio.com -> http://127.0.0.1:3500.
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.assetflow.web</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>__PROJECT_DIR__/scripts/run-web-server.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>__PROJECT_DIR__</string>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>ThrottleInterval</key>
  <integer>10</integer>

  <key>StandardOutPath</key>
  <string>__PROJECT_DIR__/data/web-server.log</string>
  <key>StandardErrorPath</key>
  <string>__PROJECT_DIR__/data/web-server.log</string>
</dict>
</plist>
