# scripts/launchd/com.assetflow.fetch-prices-us.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  launchd job: post-US-close price refresh once per day at 05:40 KST.

  Why this exists:
    The main com.assetflow.fetch-prices job runs at 15:35 KST — *before* the
    US session for that KST day opens (US opens 22:30 KST). So US symbols
    always carried the previous session's close until the next day's 15:35
    run (e.g. at 08:00 KST the app showed Friday's GOOGL close, not Monday's).

  What it does:
    Runs `fetch-prices.py --us-only` — re-fetches only NASDAQ/NYSE closes
    and merges them into the existing prices.json, leaving KRX / crypto / fx
    untouched. Sends a silent (data-only) FCM push so devices background-
    refresh without showing a dawn notification.

  Why 05:40 and not 05:00?
    US regular session closes 16:00 ET = 05:00 KST. Yahoo needs a few minutes
    to finalise the daily bar — same rationale as 15:35 vs 15:30 for KRX.

  Install / control:
    bash scripts/launchd/install.sh
    launchctl list | grep com.assetflow
    launchctl unload ~/Library/LaunchAgents/com.assetflow.fetch-prices-us.plist
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.assetflow.fetch-prices-us</string>

  <key>ProgramArguments</key>
  <array>
    <string>__PROJECT_DIR__/.venv/bin/python</string>
    <string>__PROJECT_DIR__/scripts/fetch-prices.py</string>
    <string>--us-only</string>
  </array>

  <key>WorkingDirectory</key>
  <string>__PROJECT_DIR__</string>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>5</integer>
    <key>Minute</key>
    <integer>40</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>__PROJECT_DIR__/data/fetch-prices-us.log</string>
  <key>StandardErrorPath</key>
  <string>__PROJECT_DIR__/data/fetch-prices-us.log</string>

  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
